# slapd.conf
[program:slapd]
command=/usr/local/sbin/slapd-run.sh
autorestart=true
