# root/etc/supervisor/conf.d/slapd.conf
[program:slapd]
command=/usr/local/sbin/slapd-run
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
